# crontab
*/5 0-10,16-23 * * * root node ~/Documents/hue-manager/fluxScenes.js
1 7 */1 * * root node ~/Documents/hue-manager/bingScenes.js